# launchd/com.worklogbridge.lorolog.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.worklogbridge.lorolog.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__PROJECT_DIR__/mac_collect_lorotopik_worklog.py</string>
        <string>--config</string>
        <string>__CONFIG_PATH__</string>
        <string>--mode</string>
        <string>daily</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>WORKLOGBRIDGE_CONFIG</key>
        <string>__CONFIG_PATH__</string>
        <key>TZ</key>
        <string>Asia/Seoul</string>
        <key>PYTHONPYCACHEPREFIX</key>
        <string>__HOME__/Library/Caches/WorklogBridge/python</string>
    </dict>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/lorolog-daily.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/lorolog-daily.error.log</string>
</dict>
</plist>
